# 02-packer/files/litellm.service
[Unit]
Description=LiteLLM Proxy Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openclaw
Group=openclaw
Environment="HOME=/home/openclaw"
ExecStart=/opt/litellm-venv/bin/litellm \
  --config /opt/openclaw/litellm-config.yaml \
  --port 4000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
